# ovc5/scripts/device_scripts/ethernet_utils/init_files/usb-gadget.target
[Unit]
Description=Hardware activated USB gadget
Documentation=man:systemd.special(7)
